# launchd/com.lowvram.labcore.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.lowvram.labcore</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>bash</string>
      <string>-lc</string>
      <string>test -f /Users/frido_mac/.config/low-vram-lab/env.sh &amp;&amp; source /Users/frido_mac/.config/low-vram-lab/env.sh; cd /Users/frido_mac/Projects/low_vram_institute &amp;&amp; source .venv/bin/activate &amp;&amp; lab daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/frido_mac/Projects/low_vram_institute/logs/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/frido_mac/Projects/low_vram_institute/logs/launchd.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/frido_mac/Projects/low_vram_institute</string>
  </dict>
</plist>
